# systemd/field-monitor-report.service
[Unit]
Description=field-monitor collect from servers + generate report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=%h
ExecStart=/home/ni9aii/code/field-monitor/scripts/collect-and-report.sh
StandardOutput=append:%h/.local/share/field-monitor/report.log
StandardError=append:%h/.local/share/field-monitor/report.log

[Install]
WantedBy=timers.target