# systemd/field-monitor-healthcheck.timer
[Unit]
Description=Run field-monitor healthcheck every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=field-monitor-healthcheck.service

[Install]
WantedBy=timers.target
