# systemd/field-monitor-healthcheck.service
[Unit]
Description=field-monitor healthcheck (apple-family reachability across vantage points)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash %h/code/field-monitor/scripts/field-monitor-healthcheck.sh
# script writes its own history to ~/.local/share/field-monitor/healthcheck.log
# and prints a short report to stdout (captured by journald)
