# systemd/field-monitor-report.timer
[Unit]
Description=field-monitor generate markdown report daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target